# c19vaxfinder.crontab
# m h  dom mon dow   command
# run every 4 hours
* */4 * * *     /home/c19vaxfinder/c19vaxfinder/c19vaxfinder.sh > /dev/null 2>&1
