# walgreens_co.crontab
# m h  dom mon dow   command
# run every 4 hours
#0 */4 * * *     /home/c19vaxfinder/c19vaxfinder/c19vaxfinder.sh > /dev/null 2>&1
# run every day at 7am local
*/30 * * * *     /home/c19vaxfinder/c19vaxfinder/walgreens_co.sh > /dev/null 2>&1

